# config/autostart/0-autorandr.desktop
[Desktop Entry]
Exec=autorandr horizontal
Name=autorandr
Type=Application
