# config/autostart/1-changewallpaper.desktop
[Desktop Entry]
Exec=fish -c "$DOTFILES/bin/changewallpaper"
Name=changewallpaper
Type=Application
