# config/autostart/99-dunst.desktop
[Desktop Entry]
Exec=/usr/bin/env fish -c 'dunst-autoicons'
Name=dunst
Comment=A customizable and lightweight notification-daemon
Type=Application
