# config/autostart/99-kdeconnect-indicator.desktop
[Desktop Entry]
Exec=/usr/bin/kdeconnect-indicator
Name=kdeconnect-indicator
Type=Application
