# config/autostart/99-polybar.desktop
[Desktop Entry]
Exec=fish -c "sleep 3 && $DOTFILES/bin/polybar-multimonitor"
Name=polybar
Type=Application
