# config/autostart/99-xobvolume.desktop
[Desktop Entry]
Exec=fish -c "mkfifo /tmp/xobvolume; tail -f /tmp/xobvolume | xob"
Name=xob
Type=Application
