# config/autostart/99-xobbrightness.desktop
[Desktop Entry]
Exec=fish -c "mkfifo /tmp/xobbrightness; tail -f /tmp/xobbrightness | xob"
Name=xob
Type=Application
