# config/autostart/kdeconnectd.desktop
[Desktop Entry]
Exec=/usr/lib/kdeconnectd
Name=kdeconnectd
Type=Application
