# config/autostart/polybar.desktop
[Desktop Entry]
Exec=polybar -r top
Name=startpolybar
Type=Application
