# config/autostart/sxhkd.desktop
[Desktop Entry]
Exec=sxhkd
Name=sxhkd
Type=Application
