# config/autostart/xss-lock.desktop
[Desktop Entry]
Exec=fish -c "xss-lock $DOTFILES/bin/lockscreen"
Name=xss-lock
Type=Application
